# dark-factory.service
[Unit]
Description=Pinehaven Dark Factory Orchestrator
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=gilberto
WorkingDirectory=/home/gilberto
ExecStart=/home/gilberto/.local/bin/dark-factory-orchestrator
Restart=on-failure
RestartSec=10

# Secrets (API keys) — never in git
EnvironmentFile=/home/gilberto/.config/dark-factory/secrets.env

# Logs go to a path without spaces
StandardOutput=append:/home/gilberto/.local/share/dark-factory/service.log
StandardError=append:/home/gilberto/.local/share/dark-factory/service-error.log

[Install]
WantedBy=multi-user.target
