# dark-factory.service
[Unit]
Description=Pinehaven Dark Factory Orchestrator
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=gilberto
WorkingDirectory=/home/gilberto/Desktop/Pinehaven Ventures/dark-factory
ExecStart=/home/gilberto/Desktop/Pinehaven Ventures/dark-factory/scripts/orchestrator.sh
Restart=on-failure
RestartSec=10

# Inject secrets from environment file (never commit this file)
EnvironmentFile=/home/gilberto/.config/dark-factory/secrets.env

# Logging
StandardOutput=append:/home/gilberto/Desktop/Pinehaven Ventures/dark-factory/logs/service.log
StandardError=append:/home/gilberto/Desktop/Pinehaven Ventures/dark-factory/logs/service-error.log

[Install]
WantedBy=multi-user.target
